# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:backend]
command=node /app/backend/dist/server.js
user=appuser
directory=/app/backend
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",PORT="4000",DATABASE_URL="%(ENV_DATABASE_URL)s",JWT_SECRET="%(ENV_JWT_SECRET)s",BCRYPT_ROUNDS="%(ENV_BCRYPT_ROUNDS)s",CORS_ORIGINS="%(ENV_CORS_ORIGINS)s",ADMIN_EMAIL="%(ENV_ADMIN_EMAIL)s",ADMIN_USERNAME="%(ENV_ADMIN_USERNAME)s",ADMIN_PASSWORD="%(ENV_ADMIN_PASSWORD)s",CLAUDE_CODE_MODEL="%(ENV_CLAUDE_CODE_MODEL)s",ANTHROPIC_API_KEY="%(ENV_ANTHROPIC_API_KEY)s",ANTHROPIC_BASE_URL="%(ENV_ANTHROPIC_BASE_URL)s",LLM_PROVIDER="%(ENV_LLM_PROVIDER)s",OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s",GOOGLE_GENAI_API_KEY="%(ENV_GOOGLE_GENAI_API_KEY)s",CONFLUENCE_PERSONAL_TOKEN="%(ENV_CONFLUENCE_PERSONAL_TOKEN)s",JIRA_PERSONAL_TOKEN="%(ENV_JIRA_PERSONAL_TOKEN)s",CONFLUENCE_URL="%(ENV_CONFLUENCE_URL)s",JIRA_URL="%(ENV_JIRA_URL)s",JIRA_EMAIL="%(ENV_JIRA_EMAIL)s",GITHUB_TOKEN="%(ENV_GITHUB_TOKEN)s",GITHUB_ORG="%(ENV_GITHUB_ORG)s",GIT_AUTHOR_NAME="%(ENV_GIT_AUTHOR_NAME)s",GIT_AUTHOR_EMAIL="%(ENV_GIT_AUTHOR_EMAIL)s"

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
